# deploy/aurora-bot.service
[Unit]
Description=Aurora forecast Discord bot (NOAA Ship Fairweather)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aurora
WorkingDirectory=/opt/aurora-bot
# bot.py loads .env from WorkingDirectory via python-dotenv (token, channel, SHIP_SCRIPT_PATH)
ExecStart=/opt/aurora-bot/.venv/bin/python bot.py
Restart=always
RestartSec=10
# light hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
